# scripts/serial-bridge/serial-bridge.service
# Unidad systemd para el serial bridge Arduino → n8n
#
# Instalación:
#   sudo cp serial-bridge.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable serial-bridge
#   sudo systemctl start serial-bridge
#
# Verificar logs:
#   journalctl -u serial-bridge -f
#
# Requisito: el usuario dylan debe pertenecer al grupo dialout
#   sudo usermod -aG dialout dylan
#   (reiniciar sesión para que tome efecto)

[Unit]
Description=Serial Bridge - Arduino a n8n webhook
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=dylan
Group=dialout
EnvironmentFile=/home/dylan/herramientas-ia/projects/tesis/scripts/serial-bridge/.env
ExecStart=/usr/bin/python3 /home/dylan/herramientas-ia/projects/tesis/scripts/serial-bridge/serial_bridge.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
